# crontab
0 10 * * * /usr/bin/php /root/chgsd/index.php  >/dev/null 2>&1
